# systemd/autosnap.service
[Unit]
Description=Proxmox Autosnapshot

[Service]
Type=oneshot
ExecStart=/opt/autosnap/autosnap.sh
