# systemd/cleanup.service
[Unit]
Description=Proxmox Autosnapshot Cleanup

[Service]
Type=oneshot
ExecStart=/opt/autosnap/autosnap.sh -clean
